# scripts/podman/openclaw.container.in
# OpenHearth gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into openhearth's ~/.config/containers/systemd/
# {{OPENHEARTH_HOME}} is replaced at install time.

[Unit]
Description=OpenHearth gateway (rootless Podman)

[Container]
Image=openhearth:local
ContainerName=openhearth
UserNS=keep-id
Volume={{OPENHEARTH_HOME}}/.openhearth:/home/node/.openhearth
EnvironmentFile={{OPENHEARTH_HOME}}/.openhearth/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
